# etc/hotplug.d/iface/65-log-hostname
#!/bin/sh

[ "$INTERFACE" = "wan" ] && {
	logger -t $(cat /proc/sys/kernel/hostname) hostname $ACTION
}
