# etc/crontabs/root
*/5 * * * *	/usr/bin/client-status.sh -t > /tmp/users.html
*/5 * * * *	sleep $(cat /tmp/rand); /usr/bin/client-status.sh -f
0 * * * *	logger -t $(cat /proc/sys/kernel/hostname) hostname
0 5 * * *	/etc/init.d/olsrd restart
