# etc/crontabs/root
*/5 * * * *	/usr/bin/client-status.sh -t > /tmp/users.html
*/5 * * * *	/usr/bin/client-status.sh -f
