# etc/crontabs/root
*/5 * * * *	/usr/bin/client-status.sh -t > /tmp/users.html
0 * * * *	logger -t $(cat /proc/sys/kernel/hostname) hostname
*/7 * * * *	if ! pidof openvpn ; then echo openvpn not running ; /etc/init.d/openvpn restart ; fi
*/5 * * * *	curl -m 10 -s "https://node:zRJOvia9DuIw@personaltelco.net/dns/dnsupdate.php?hostname=$(cat /proc/sys/kernel/hostname)" > /dev/null 2>&1
*/5 * * * *	sleep $(expr $(cat /tmp/rand) + 1) ; host-ipaddr-client knod.ptp -4
*/5 * * * *	sleep $(expr $(cat /tmp/rand) + 2) ; host-ipaddr-client knod.personaltelco.net -4
*/5 * * * *	sleep $(expr $(cat /tmp/rand) + 3) ; host-ipaddr-client knod.personaltelco.net -6
