# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/etc/init.d/snmpd restart
dd if=/dev/urandom bs=4 count=1 2>/dev/null | md5sum | tr -d 'a-z' | awk '{ v = substr($1,0,9) ; srand(v) ; print int(295 * rand()) }' > /tmp/rand
exit 0
