# etc/uci-defaults/ptp.private.defaults
#!/bin/sh

uci batch <<EOF
set network.pub.ifname="PTP_PUBIFACES_PTP"
set network.priv=interface
set network.priv.ifname="PTP_PRIVIFACES_PTP"
set network.priv.type=bridge
set network.priv.proto=static
set network.priv.ipaddr=PTP_PRIVADDR_PTP
set network.priv.netmask=PTP_PRIVNETMASK_PTP
commit network

set dhcp.priv=dhcp
set dhcp.priv.interface=priv
set dhcp.priv.start=100
set dhcp.priv.limit=150
set dhcp.priv.leasetime=12h
set dhcp.priv.force=1
commit dhcp
EOF
